# scripts/com.claude-code.autostart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-code.autostart</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>$HOME/.claude/scripts/start-all.sh</string>
    <string>--silent</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>$HOME/.claude/logs/autostart.log</string>

  <key>StandardErrorPath</key>
  <string>$HOME/.claude/logs/autostart-error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>$HOME</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
